# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null

[program:uvicorn]
command=uvicorn main:app --host 0.0.0.0 --port 38888
directory=/app
autostart=true
autorestart=true
stopsignal=TERM
stopasgroup=true
killasgroup=true

stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stopsignal=QUIT

stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0